# motorclaim.service
[Unit]
Description=Motor Claim Decision API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/Motorclaimdecision_main
Environment="MOTORCLAIM_BASE_DIR=/opt/Motorclaimdecision_main"
Environment="PYTHONPATH=/opt/Motorclaimdecision_main"
Environment="POPPLER_PATH=/opt/Library/poppler-25.07.0/bin"
Environment="PATH=/opt/Library/poppler-25.07.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/usr/bin/python3 unified_api_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
